# docker/s6/s6-rc.d/web/run
#!/command/with-contenv bash
set -euo pipefail

mode="${APP_ENV:-prod}"

if [[ "$mode" == "dev" ]]; then
  echo "[web] dev mode: vite dev"
  cd /app/web
  # bind on 0.0.0.0 so Docker can hit it from host/other containers
  exec npm run dev -- --host 0.0.0.0 --port 5173
else
  echo "[web] prod mode: static assets already built at /usr/share/nginx/html, served by nginx service"
  # In prod, this service is not needed but must stay running for s6
  exec sleep infinity
fi
