# docker/s6/s6-rc.d/postgres-ready/up
bash /etc/s6-overlay/scripts/wait-for-postgres.sh
