# docker/s6/s6-rc.d/prowlarr/run
#!/usr/bin/env bash
set -euo pipefail

PROWLARR_PORT="9696"
PROWLARR_API_KEY="prowlarr-api-key"

# Export the API key for the current process
export PROWLARR__SERVER__PORT="$PROWLARR_PORT"
export PROWLARR__AUTH__APIKEY="$PROWLARR_API_KEY"

# DB config
export PROWLARR__POSTGRES__HOST="localhost"
export PROWLARR__POSTGRES__PORT="5432"
export PROWLARR__POSTGRES__USER="prowlarr"
export PROWLARR__POSTGRES__PASSWORD="prowlarrpw"

# Start Prowlarr with the API key
exec /opt/prowlarr/Prowlarr -nobrowser -data=/var/lib/prowlarr/ "$@"
