# docker/s6/s6-rc.d/postgres/run
#!/usr/bin/env bash
set -euo pipefail

: "${PGDATA:=/var/lib/postgresql/data}"
: "${PG_BINDIR:=/usr/lib/postgresql/16/bin}"

export PATH="$PG_BINDIR:$PATH"

echo "[postgres] starting postgres with PGDATA=${PGDATA}"
exec s6-setuidgid postgres postgres -D "$PGDATA" -c listen_addresses=*
