# docker/s6/s6-rc.d/postgres-init/up
bash /etc/s6-overlay/scripts/postgres-init.sh
